# autostart/ukui-panel.desktop.in
[Desktop Entry]
Type=Application
Name=Panel
TryExec=ukui-panel
Exec=ukui-panel
OnlyShowIn=UKUI;
#TRANSLATIONS_DIR=translations
